# certbot.service
[Unit]
Description=Certbot Script
[Service]
Type=oneshot
ExecStart=/opt/certbot.sh
PrivateTmp=true